# supervisor-xvfb.conf
[program:xvfb]
directory=/tmp
command=Xvfb
environment=DISPLAY=":0"
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/xvfb.err.log
stdout_logfile=/var/log/xvfb.out.log
